# gentoo-init.d/autothrottle
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="autothrottle daemon"
description="daemon that reduces the maximum cpu frequency to avoid overheating by d3da"
command=autothrottle
command_args="${autothrottle_args}"

